# systemd/pia-wg-monitor.service
[Unit]
Description=PIA WireGuard Connection Monitor
After=pia-wg.service
Requires=pia-wg.service

[Service]
Type=oneshot
ExecStartPre=/bin/sleep 30
ExecStart=/usr/local/bin/pia-monitor.sh
TimeoutStartSec=120
User=root
Group=root

[Install]
WantedBy=multi-user.target
